# docker/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ─── Web Server (Sail 8.4 runtime) ───────────────────────────────────────────
[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ─── Queue Worker ─────────────────────────────────────────────────────────────
[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=1 --timeout=3600 --max-jobs=100 --queue=default
process_name=%(program_name)s_%(process_num)02d
user=sail
numprocs=1
autostart=true
autorestart=true
startsecs=1
startretries=3
priority=10
stdout_logfile=/tmp/queue-worker.log
stderr_logfile=/tmp/queue-worker-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB